# install_scripts/SYSTEM/etc/systemd/system/altimeter.service
[Unit]
Description=Altimeter Gui Service

[Service]
ExecStart=/home/pi/start-altimeter-gui
Restart=always

[Install]
WantedBy=multi-user.target
