# install_scripts/SYSTEM/etc/systemd/system/ser-mon.service
[Unit]
Description=SerialAltimeterMonitor

[Service]
ExecStart=/home/pi/start-serial-monitor.sh
Restart=always

[Install]
WantedBy=multi-user.target
